# lib/dpkg/info/postgresql-common.conffiles
/etc/init.d/postgresql
/etc/logrotate.d/postgresql-common
/etc/postgresql-common/pg_upgradecluster.d/analyze
